# configs/butane/files/harbor-install.service
[Unit]
Description=Install and start Harbor
After=network-online.target docker.service
Requires=docker.service
ConditionPathExists=!/var/harbor/install/harbor/docker-compose.yml

[Service]
Type=oneshot
RemainAfterExit=true
Environment=HARBOR_VERSION=v2.14.2
ExecStart=/usr/bin/bash /var/harbor/install/setup-harbor.sh
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target